# backsource/config/systemd/user/wal-swap.timer
[Unit]
Description=Ejecutar wal-swap cada 20 minutos

[Timer]
OnBootSec=1min
OnUnitActiveSec=20min
Persistent=true

[Install]
WantedBy=timers.target
